# crontab
SHELL=/bin/bash

*/15 * * * * root /opt/FullView/bot_cronjob.sh &>> /opt/FullView/bot_cronjob.log
35 3,11,19 * * * root /opt/FullView/bgpdump-stat.sh &>> /opt/FullView/bot_cronjob.log